# build/config/supervisor.conf
[program:undockd]
command = uwsgi --ini /config/undockd.conf

[program:nginx]
command = service nginx start
