# build/config/supervisor.conf
[program:undockd]
command = uwsgi --ini /config/undockd.ini

[program:nginx]
command = service nginx restart

[program:webpack]
directory = /application/undockd/static
command = webpack --watch
