# build/config/supervisor.ini
[program:undockd]
command = uwsgi --ini /config/undockd.ini

[program:nginx]
command = service nginx restart

[program:tsc]
directory = /application/undockd/static
command = webpack --watch
